# docker/cron.d/cronjob
*/2 * * * * /usr/local/bin/php /var/www/html/cli/sync_site_products.php >> /var/log/cron.log 2>&1
0 * * * * /usr/local/bin/php /var/www/html/cli/clear_cron_log.php >> /var/log/cron.log 2>&1
